# templates/default/node-app-init.erb
#!/bin/bash
################################################
# <%= @node_application %>
#
# script used to start and stop the node <%= @node_application %> service
#
###############################################
#
# chkconfig: 345 90 16
# description: <%= @node_application %> node application

. /etc/rc.d/init.d/functions

USER="<%= @user %>"
NODE_APP="<%= @node_application %>"
DAEMON="<%= @start_application %>"
ROOT_DIR="<%= @app_root_directory %>"
START_COMMAND="<%= @start_command %>"
STOP_COMMAND="<%= @stop_command %>"

LOCK_FILE="/var/lock/subsys/<%= @node_application %>"

do_start()
{
        if [ ! -f "$LOCK_FILE" ]; then
                cd $ROOT_DIR
                runuser -l "$USER" -c "cd $ROOT_DIR; $DAEMON $START_COMMAND $NOD_APP"
                RETVAL=$?
                echo
                [ $RETVAL -eq 0 ] && touch $LOCK_FILE
        else
                echo "$NOD_APP is locked"
        fi
}

do_stop()
{
        echo -n $"Stopping $NODE_APP"
        runuser -l "$USER" -c "cd $ROOT_DIR; $DAEMON $STOP_COMMAND $NOD_APP"
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
}

case "$1" in
        start)
                do_start
                ;;
        stop)
                do_stop
                ;;
        restart)
                do_stop
                do_start
                ;;
        *)
                echo "Usage: $0 {start|stop|restart}"
                RETVAL=1
esac
exit $RETVAL
